# rclone/rootfs/etc/services.d/rclone/run
#!/usr/bin/with-contenv bashio

RCLONE_CONFIG=$(bashio::config 'configuration_path')
REMOTES=$(bashio::config 'remotes')

for remote in ${REMOTES}; do
    remote_name=$(bashio::jq "${remote}" ".name")
    remote_path=$(bashio::jq "${remote}" ".path")
    remote_retention=$(bashio::jq "${remote}" ".retention_days")

    bashio::log.info "Starting remote copy for ${remote_name}..."
    
    rclone --config ${RCLONE_CONFIG} \
        sync -v \
        --max-age "${remote_retention}d" \
        --delete-during \
        --ignore-errors \
        --ignore-existing "${remote_name}:${remote_path}" /share/rd/
    
    bashio::log.info "Remote copy for ${remote_name} finished"
done
